# deploy/kaix-backend.service
[Unit]
Description=KaiX unified backend (Python ThreadingHTTPServer)
After=network.target
Wants=network-online.target
# If the process keeps crashing, give up after 5 attempts in 60s so
# systemd doesn't hammer a broken binary forever. journalctl will
# show the failure cause.
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=kaix
Group=kaix
WorkingDirectory=/opt/kaix/web
EnvironmentFile=/etc/kaix.env
# Primary instance: listens on 8787 and owns the singleton crawler scheduler.
# Extra capacity is added via kaix-backend-worker@<port>.service (8788+),
# which set KAIX_ENABLE_SCHEDULERS=0. All instances share the SQLite WAL DB.
Environment=KAIX_PORT=8787
Environment=KAIX_ENABLE_SCHEDULERS=1
ExecStart=/usr/bin/python3 -u /opt/kaix/web/server.py
Restart=always
RestartSec=3
# graceful: server.py traps SIGTERM and runs server.shutdown()
KillSignal=SIGTERM
TimeoutStopSec=15

# hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/opt/kaix/web
# Each connection takes ~2 fds; raise the cap so we can serve many
# concurrent users + SSE streams without EMFILE.
LimitNOFILE=65536
# Cap memory so a runaway request can't OOM the host.
MemoryHigh=1500M
MemoryMax=2048M
# Reasonable file size cap (1GB) — prevents an upload bug from filling disk.
LimitFSIZE=1073741824

# tame log output (journalctl handles rotation)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kaix-backend

[Install]
WantedBy=multi-user.target
